# docker/supervisord.local.conf
[supervisord]
nodaemon=false

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:php-fpm]
command=/usr/sbin/php-fpm8.2 --nodaemonize
autostart=true
autorestart=true

[program:mariadb]
command=/usr/bin/mysqld_safe --datadir=/var/lib/mysql
autostart=true
autorestart=true

[program:redis]
command=/usr/bin/redis-server
autostart=true
autorestart=true